# systemd/jinyang-backlog.service
[Unit]
Description=Backlog Processor - Session recovery queue runner
After=network.target

[Service]
Type=oneshot
User=%u
Group=%u
ExecStart=%h/Applications/jinyang/scripts/backlog-processor.sh
Nice=10
IOSchedulingClass=idle
IOSchedulingPriority=7
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=jinyang-backlog

[Install]
WantedBy=multi-user.target
